# environment/.config/profile.d/00-lang.sh
#!/usr/bin/env sh

## language settings

export LC_ALL="${LC_ALL:-en_US.UTF-8}"
export LANG="${LANG:-en_US.UTF-8}"
