# environment/.config/profile.d/50-defaults.sh
#!/usr/bin/env sh

## default applications

export EDITOR="nvim"
export VISUAL="nvim"
export TERMINAL="kitty"
export BROWSER="librewolf"
export BROWSER_PRIVATE_OPT="--private-window"
export SYSTEMD_EDITOR="$EDITOR"
export READER="zathura"
export VIDEO="mpv"
export IMAGE="imv"
export FILES="nemo"
